# frontend_v3/html/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:npm]
command=npm run dev
directory=/app
autostart=true
autorestart=true
